# deploy/operator-os.service
[Unit]
Description=WE302 Operator OS (Western Electric 302 exchange)
After=multi-user.target sound.target
Wants=network-online.target

[Service]
Type=simple
User=john
Group=john
WorkingDirectory=/home/john/operator
Environment=PATH=/home/john/.local/bin:/usr/bin
# Prefer uv so the project venv + system-site-packages gpio stack stay consistent.
ExecStart=/home/john/.local/bin/uv run operator-os run --config config/hardware_profile.yaml
Restart=on-failure
RestartSec=3
# Fail safe: kill leftover aplay/arecord on stop.
KillMode=control-group
TimeoutStopSec=10

[Install]
WantedBy=multi-user.target
